# etc/rc.d/adjkerntz
#!/bin/sh
#
# $FreeBSD: release/7.0.0/etc/rc.d/adjkerntz 168283 2007-04-02 22:53:07Z des $
#

# PROVIDE: adjkerntz
# REQUIRE: FILESYSTEMS random
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr

name="adjkerntz"
start_cmd="adjkerntz -i"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
