# etc/rc.d/adjkerntz
#!/bin/sh
#
# $FreeBSD: release/10.0.0/etc/rc.d/adjkerntz 239570 2012-08-22 18:49:02Z obrien $
#

# PROVIDE: adjkerntz
# REQUIRE: FILESYSTEMS postrandom
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr

name="adjkerntz"
start_cmd="adjkerntz -i"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
