# etc/rc.d/rarpd
#!/bin/sh
#
# $FreeBSD: stable/9/etc/rc.d/rarpd 242153 2012-10-26 18:06:49Z obrien $
#

# PROVIDE: rarpd
# REQUIRE: DAEMON FILESYSTEMS
# BEFORE:  LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="rarpd"
rcvar="rarpd_enable"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/ethers"

load_rc_config $name
run_rc_command "$1"
