# etc/rc.d/rpcbind
#!/bin/sh
#
# $FreeBSD: stable/9/etc/rc.d/rpcbind 231653 2012-02-14 10:16:56Z dougb $
#

# PROVIDE: rpcbind
# REQUIRE: NETWORKING ntpdate syslogd named
# KEYWORD: shutdown

. /etc/rc.subr

name="rpcbind"
rcvar="rpcbind_enable"
command="/usr/sbin/${name}"

stop_postcmd='/bin/rm -f /var/run/rpcbind.*'

load_rc_config $name
run_rc_command "$1"
