# etc/rc.d/rwho
#!/bin/sh
#
# $FreeBSD: stable/9/etc/rc.d/rwho 231653 2012-02-14 10:16:56Z dougb $
#

# PROVIDE: rwho
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="rwhod"
rcvar="rwhod_enable"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
