# etc/rc.d/rwho
#!/bin/sh
#
# $FreeBSD: release/7.0.0/etc/rc.d/rwho 175737 2008-01-28 07:58:32Z dougb $
#

# PROVIDE: rwho
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="rwhod"
rcvar="`set_rcvar`"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
