# etc/rc.d/timed
#!/bin/sh
#
# $FreeBSD: release/10.0.0/etc/rc.d/timed 230099 2012-01-14 02:18:41Z dougb $
#

# PROVIDE: timed
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="timed"
rcvar="timed_enable"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
