# etc/rc.d/timed
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/timed,v 1.5 2004/10/07 13:55:26 mtm Exp $
# $MidnightBSD: src/etc/rc.d/timed,v 1.2 2006/08/28 02:28:59 laffer1 Exp $

# PROVIDE: timed
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="timed"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
