# etc/periodic/monthly/Makefile
# $FreeBSD: release/7.0.0/etc/periodic/monthly/Makefile 100316 2002-07-18 12:33:01Z ru $

FILES=	200.accounting \
	999.local

.include <bsd.prog.mk>
